# deploy/trading-bot.service
[Unit]
Description=Safe Yield Trading Bot
After=network.target

[Service]
Type=simple
User=YOUR_NUC_USERNAME
WorkingDirectory=/home/YOUR_NUC_USERNAME/safe-yield-stock-bot
ExecStart=/home/YOUR_NUC_USERNAME/safe-yield-stock-bot/.venv/bin/python run_telegram.py
Restart=on-failure
RestartSec=10
EnvironmentFile=/home/YOUR_NUC_USERNAME/safe-yield-stock-bot/.env

[Install]
WantedBy=multi-user.target
